# .crontab
# m h dom mon dow command
* * * * * for i in $(seq 5 15 59); do (sleep $i; /usr/bin/cron/sorapiyo-notify.sh) & done
