# .crontab
# m h dom mon dow command
* * * * * date >> /var/tmp/piyo
